# mkosi.skeleton/etc/profile.d/environment.sh
export DISPLAY=:0
systemctl --user import-environment DISPLAY
